# projects/zFPGA_DE10_NANO_SoC/12_linux_app/scripts/auto_run.sh
#!/bin/sh
#chkconfig: 2345 80 90
#description:开机自动启动的脚本程序
cd /home/root
/home/root/my_first_hps-fpga r